# roles/hermes/templates/hermes-dashboard.service.j2
[Unit]
Description=Hermes Web Dashboard
Documentation=https://hermes-agent.nousresearch.com/docs/user-guide/features/web-dashboard
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User={{ hermes_user }}
Group={{ hermes_group }}
WorkingDirectory={{ hermes_home_dir }}
Environment=HOME={{ hermes_home_dir }}
Environment=HERMES_HOME={{ hermes_home_dir }}/.hermes-{{ hermes_dashboard_instance }}
EnvironmentFile={{ hermes_home_dir }}/.hermes-{{ hermes_dashboard_instance }}/.env
ExecStart={{ hermes_venv_dir }}/bin/hermes dashboard --host {{ hermes_dashboard_host }} --port {{ hermes_dashboard_port }}{{ ' --no-open' if hermes_dashboard_no_open | bool else '' }}
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=hermes-dashboard

[Install]
WantedBy=multi-user.target
