# roles/hermes/templates/hermes-gateway.service.j2
[Unit]
Description=Hermes Messaging Gateway ({{ hermes_instance.name }})
Documentation=https://hermes-agent.nousresearch.com/docs/
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User={{ hermes_user }}
Group={{ hermes_group }}
WorkingDirectory={{ hermes_home_dir }}
Environment=HOME={{ hermes_home_dir }}
Environment=HERMES_HOME={{ hermes_config_dir }}/profiles/{{ hermes_instance.name }}
EnvironmentFile={{ hermes_config_dir }}/profiles/{{ hermes_instance.name }}/.env
ExecStart={{ hermes_venv_dir }}/bin/hermes gateway run
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=hermes-gateway-{{ hermes_instance.name }}

[Install]
WantedBy=multi-user.target
