# roles/litellm/templates/litellm.service.j2
[Unit]
Description=LiteLLM AI Model Proxy
After=docker.service
Requires=docker.service

[Service]
Type=exec
WorkingDirectory={{ litellm_project_dir }}
ExecStart=/usr/bin/docker compose -f {{ litellm_compose_file }} up --remove-orphans
ExecStop=/usr/bin/docker compose -f {{ litellm_compose_file }} down
Restart=on-failure
RestartSec=10s
User={{ litellm_primary_user }}
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
